# .config/autostart/idlehandler.wayland.desktop
[Desktop Entry]
Exec=~/bin/idlehandler.wayland
Name=swayidle
Type=Application
OnlyShowIn=sway;
