# .config/autostart/idlehandler.wayland.desktop
[Desktop Entry]
Exec=idlehandler.wayland
Name=swayidle
Type=Application
OnlyShowIn=sway;
