# .config/autostart/idlehandler.wayland.desktop
[Desktop Entry]
Exec=idlehandler.wayland
#Exec=systemctl --user start idlehandler
Name=swayidle
Type=Application
OnlyShowIn=sway;
