# .config/autostart/kanshi.desktop
[Desktop Entry]
Exec=/usr/bin/kanshi
Name=kanshi
Comment=Dynamic output configuration
Type=Application
OnlyShowIn=sway;
