# .config/autostart/swayrd.desktop
[Desktop Entry]
#Hidden=true
Exec=/home/cel/bin/swayrd
Name=swayrd
Comment=swayrd deamon for sway
Type=Application
OnlyShowIn=sway;
