# .config/autostart/tracker-miner-rss.desktop
[Desktop Entry]
Exec=/usr/libexec/tracker-miner-rss
Name=tracker-miner-rss
Type=Application
Hidden=true
