# .config/autostart/tracker-miner-fs.desktop
[Desktop Entry]
Exec=/usr/libexec/tracker-miner-fs
Name=tracker-miner-fs
Type=Application
Hidden=true
