# .config/autostart/tracker-extract.desktop
[Desktop Entry]
Exec=/usr/libexec/tracker-extract
Name=tracker-extract
Type=Application
Hidden=true
