# .config/autostart/tracker-miner-apps.desktop
[Desktop Entry]
Exec=/usr/libexec/tracker-miner-apps
Name=tracker-miner-apps
Type=Application
Hidden=true
